# templates/mariabackup.timer.j2
[Unit]
Description=MariaDB Backup timer

[Timer]
OnCalendar={{ cron }}
Persistent=true

[Install]
WantedBy=timers.target
