# templates/mariabackup.service.j2
[Unit]
Description=MariaDB Backup service

[Service]
Environment="BACKUP_ROOT={{ backup_root}}"
ExecStart={{ scripts_dir }}/backup.sh
